# whisper-dictate.service
[Unit]
Description=Whisper push-to-talk dictation daemon
After=graphical-session.target
ConditionPathIsDirectory=/sys/module/amdgpu

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStartPre=/bin/sh -c 'until xset q >/dev/null 2>&1; do sleep 2; done'
ExecStart=%h/Config/whisper/whisper-dictate
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
